# deploy/systemd/jobseek-crawler-reconciliation.service
[Unit]
Description=Jobseek bounded cross-store posting reconciliation
Documentation=https://github.com/colophon-group/jobseek/blob/main/docs/03-crawler-architecture.md#cross-store-reconciliation
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target
ConditionPathExists=/usr/local/sbin/jobseek-crawler-reconciliation
ConditionPathExists=/home/deploy/.env

[Service]
Type=oneshot
User=deploy
Group=deploy
ExecStart=/usr/local/sbin/jobseek-crawler-reconciliation
TimeoutStartSec=3h
TimeoutStopSec=90s
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
UMask=0077

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=read-only
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
ReadWritePaths=/run/lock

[Install]
WantedBy=multi-user.target
