# deploy/systemd/jobseek-codex-daily-error-review.timer
# Install as /etc/systemd/system/jobseek-codex-daily-error-review.timer.
[Unit]
Description=Schedule the Jobseek Codex daily crawler error review
Documentation=file:///srv/jobseek-codex/repo/docs/18-codex-automation-deployment.md

[Timer]
OnCalendar=*-*-* 09:00:00 UTC
RandomizedDelaySec=10min
AccuracySec=1min
Persistent=false
Unit=jobseek-codex-daily-error-review.service

[Install]
WantedBy=timers.target
